# profile.d/vsync_off
# VSync off by default
export vblank_mode=0
