# autostart/hubic.desktop
[Desktop Entry]
Name=hubiC
Exec=/usr/lib/hubic/hubic-service
Name=hubiC
Icon=hubic
Type=Application
