# autostart/numlockx.desktop
[Desktop Entry]
Exec=numlockx
Name=numlockx
Terminal=false
Type=Application
