# autostart/screenstalk.desktop
[Desktop Entry]
Exec=screenstalk -t 10m
Name=screenstalk
Terminal=false
Type=Application
