# autostart/ssh-agent.desktop
[Desktop Entry]
Exec=ssh-agent
Name=ssh-agent
Terminal=false
Type=Application
